# build/s6/samba/run
#!/bin/bash

# Treat unset variables as an error
chmod +x /s6-bin/samba/entrypoint.sh

if ! [[ -z ${SAMBAENABLE} ]]; then
    # SAMBA SCRIPT START
    echo "-----------------------------------------------------------------"
    echo "Starting SAMBA/CIFS Server..."
    echo "-----------------------------------------------------------------"

    # Create arguments string
    # Enable NMBD
    ARGUMENTS="-n"
    # Server Name
    ARGUMENTS="${ARGUMENTS} -g \"netbios name = \"${SRVNAME}\"\""
    # Workgroup Name
    ARGUMENTS="${ARGUMENTS} -w \"${WORKGROUPNAME}\""
    # Parse and split Users by :
    while IFS=':' read -ra TEMP; do
        for i in "${TEMP[@]}"; do
            ARGUMENTS="${ARGUMENTS} -u \"${i}\""
        done
    done <<< "${USERS}"
    # Parse and split Mounts by :
    while IFS=':' read -ra TEMP; do
        for i in "${TEMP[@]}"; do
            ARGUMENTS="${ARGUMENTS} -s \"${i}\""
        done
    done <<< "${MOUNTS}"

    # Call the helper script
    eval /s6-bin/samba/entrypoint.sh "${ARGUMENTS}"

    while smbclient -L '\\localhost' -U '%' -m SMB3 >& /dev/null;:
    do
        echo "-----------------------------------------------------------------"
        echo "Samba server still running and healthy."
        echo "Sleeping ${SLEEPTIME:-3600} minutes."
        echo "-----------------------------------------------------------------"
        sleep ${SLEEPTIME:-3600}
    done
else
    echo "-----------------------------------------------------------------"
    echo "Samba Server not enabled."
    echo "-----------------------------------------------------------------"
    sleep 86400
fi