# build/s6/.s6-svscan/finish
#!/bin/bash

echo "--------------------------------------------------------------------------"
echo "Terminate Signal caught."
echo "--------------------------------------------------------------------------"

echo
echo "--------------------------------------------------------------------------"
echo "Calling all finish scripts."
echo "--------------------------------------------------------------------------"
echo
for f in ../**/finish; do  
    echo
    echo "--------------------------------------------------------------------------"
    echo "Calling finish script: ${f}"
    echo "--------------------------------------------------------------------------"
    echo
    bash "$f" -H
done