# build/services.d/softether-vpncli/finish
#!/bin/bash

source /scripts/logger.sh

log_interrupt "Terminating softether-vpncli..." "both"

log_warn "Disconnecting from the server..."
/s6-bin/softether-vpncli/vpncmd localhost /CLIENT /CMD accountdisconnect ${CONNNAME}

log_warn "Deleting Virtual Network Adapter..."
/s6-bin/softether-vpncli/vpncmd localhost /CLIENT /CMD nicdelete vpn${ADAPTERNAME}

log_warn "Shutting down the SoftEther daemon..."
/s6-bin/softether-vpncli/vpnclient stop

log_interrupt "Terminated softether-vpncli." "both"
