# build/s6/softether-vpncli/finish
#!/bin/bash

echo "-----------------------------------------------------------------"
echo "Terminate signal caught."
echo "Terminating VPN Connection"
echo "-----------------------------------------------------------------"
/s6-bin/softether-vpncli/vpncmd localhost /CLIENT /CMD accountdisconnect ${CONNNAME} 
echo "-----------------------------------------------------------------"
echo "Deleting Virtual Network Adapter"
echo "-----------------------------------------------------------------"
/s6-bin/softether-vpncli/vpncmd localhost /CLIENT /CMD nicdelete vpn${ADAPTERNAME}
echo "-----------------------------------------------------------------"
echo "Shutting down the daemon"
echo "-----------------------------------------------------------------"
/s6-bin/softether-vpncli/vpnclient stop
echo "-----------------------------------------------------------------"
echo "ALL DONE"
echo "-----------------------------------------------------------------"
echo "softether-vpncli ready to terminate..."
echo "-----------------------------------------------------------------"
echo 