# build/s6/softether-vpncli/finish
#!/bin/bash

echo 
echo "--------------------------------------------------------------------------"
echo "Terminate signal caught."
echo 
echo "Terminating VPN Connection"
echo "--------------------------------------------------------------------------"
echo 
/s6-bin/softether-vpncli/vpncmd localhost /CLIENT /CMD accountdisconnect ${CONNNAME} 
echo 
echo "Deleting Virtual Network Adapter"
echo "--------------------------------------------------------------------------"
echo 
/s6-bin/softether-vpncli/vpncmd localhost /CLIENT /CMD nicdelete ${ADAPTERNAME}
echo 
echo "Shutting down the daemon"
echo "--------------------------------------------------------------------------"
echo 
/s6-bin/softether-vpncli/vpnclient stop
echo 
echo "--------------------------------------------------------------------------"
echo "ALL DONE"
echo "--------------------------------------------------------------------------"
echo "softether-vpncli ready to terminate..."
echo 