# build/s6/softether-vpncli/run
#!/bin/bash

# SOFTETHER START
echo
echo "Starting Softether Client..."
echo "--------------------------------------------------------------------------"
echo 
/s6-bin/softether-vpncli/vpnclient start > /dev/null
echo
echo "Creating new virtual ethernet adapter."
echo "--------------------------------------------------------------------------"
echo 
/s6-bin/softether-vpncli/vpncmd localhost /CLIENT /CMD nicdelete ${ADAPTERNAME}  > /dev/null
/s6-bin/softether-vpncli/vpncmd localhost /CLIENT /CMD niccreate ${ADAPTERNAME} > /dev/null
echo
echo "Importing defined ${CONNNAME} connection."
echo "--------------------------------------------------------------------------"
echo
/s6-bin/softether-vpncli/vpncmd localhost /CLIENT /CMD accountdisconnect ${CONNNAME} > /dev/null
/s6-bin/softether-vpncli/vpncmd localhost /CLIENT /CMD accountdelete ${CONNNAME} > /dev/null
# Writing this with expect here because command line interface does not allow / in start of filename
# Maybe there is another way of doing that.
/usr/bin/expect > /dev/null<<EOD
set timeout -1
spawn /s6-bin/softether-vpncli/vpncmd localhost /CLIENT /CMD
expect "VPN Client>"
send "accountimport\r"
expect "Import Source File Name: "
send "/s6-bin/softether-vpncli/${CONNNAME}.vpn\r"
expect "VPN Client>"
send "exit\r"
EOD

if ! [[ -z ${MACADD} ]]; then
    echo
    echo "MAC Address variable defined in enviromental variables."
    echo "Setting MAC address of the interface to: ${MACADD}"
    echo "--------------------------------------------------------------------------"
    echo    
    /s6-bin/softether-vpncli/vpncmd localhost /CLIENT /CMD nicsetsetting ${ADAPTERNAME} /MAC:${MACADD} > /dev/null
    /s6-bin/softether-vpncli/vpncmd localhost /CLIENT /CMD nicupgrade ${ADAPTERNAME} > /dev/null
    echo
    echo "Restart VPN Client for MAC Address to be set."
    echo "--------------------------------------------------------------------------"
    echo 
    /s6-bin/softether-vpncli/vpnclient stop > /dev/null
    /s6-bin/softether-vpncli/vpnclient start > /dev/null
fi

echo
echo "Connecting to defined ${CONNNAME} connection."
echo "--------------------------------------------------------------------------"
echo
/s6-bin/softether-vpncli/vpncmd localhost /CLIENT /CMD accountconnect ${CONNNAME} > /dev/null

# This is only valid for internal connections.
if ! [[ -z ${INTCONN} ]]; then
    echo
    echo "Internal connection variable defined in enviromental variables."
    echo "Set adapter MTU to 1200."
    echo "--------------------------------------------------------------------------"
    echo
    ifconfig vpn_vpn mtu 1200 > /dev/null
    sleep 1 > /dev/null
    ifconfig vpn_vpn down && ifconfig vpn_vpn up > /dev/null
fi

echo
echo "Tailing log files..."
echo "--------------------------------------------------------------------------"
echo
tail -F /usr/local/libexec/softether/vpnclient/*_log/*.log &
echo
echo "Requesting IP Address from Server"
echo "--------------------------------------------------------------------------"
echo
dhclient vpn_vpn > /dev/null
IPADDRESS=$(ifconfig vpn_vpn | awk '/inet addr/{print substr($2,6)}')
if ! [[ -z "${IPADDRESS}" ]]; then
    echo
    echo "IP Address Obtained: ${IPADDRESS}"
    echo "--------------------------------------------------------------------------"
    echo
    NETWORKADDRESS=$(echo ${IPADDRESS} | sed -r -e 's/^([0-9]*\.[0-9]*\.[0-9]*\.).*$/\1/')
    while ping -c 1 -W 10 "${NETWORKADDRESS}1" >& /dev/null;:
    do
        echo "--------------------------------------------------------------------------"
        echo "IP Address: ${IPADDRESS}"
        echo "Connection to ${NETWORKADDRESS}1 can be established. User has valid IP configuration."
        echo "Sleeping 10 minutes."
        echo "--------------------------------------------------------------------------"
        sleep 600
    done
else
    echo "err: Can not obtain IP Address."
fi