# build/services.d/samba/run
#!/bin/bash

source /scripts/logger.sh

if ! [[ -z ${SAMBAENABLE} ]]; then
  # SAMBA SCRIPT START
  log_start "Starting SAMBA/CIFS Server..."

  # Create arguments string
  # Enable NMBD
  ARGUMENTS="-n"
  # Server Name
  ARGUMENTS="${ARGUMENTS} -g \"netbios name = \"${SRVNAME}\"\""
  # Workgroup Name
  ARGUMENTS="${ARGUMENTS} -w \"${WORKGROUPNAME}\""

  # Parse and split Users by :
  while IFS=':' read -ra TEMP; do
    for i in "${TEMP[@]}"; do
      ARGUMENTS="${ARGUMENTS} -u \"${i}\""
    done
  done <<<"${USERS}"

  # Parse and split Mounts by :
  while IFS=':' read -ra TEMP; do
    for i in "${TEMP[@]}"; do
      ARGUMENTS="${ARGUMENTS} -s \"${i}\""
    done
  done <<<"${MOUNTS}"

  # Call the helper script
  eval ./entrypoint.sh "${ARGUMENTS}"

  while
    smbclient -L '\\localhost' -U '%' -m SMB3 >&/dev/null
    :
  do
    log_this "Samba client can connect to self, server still healthy." "${GREEN}SAMBA${RESET}" "INFO"
    log_wait "Next check in ${SLEEPTIME:-3600} seconds."
    s6-sleep ${SLEEPTIME:-3600}
  done
else
  log_finish "Samba Server not enabled."
  s6-sleep 86400
fi
